# galaxyenv/files/orval/rsyncd.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
Service manifest for cswrsyncd
Autogenerated by CSWinitsmf
2008-04-21 Peter Bonivart
-->

<service_bundle type='manifest' name='rsyncd:rsyncd'>

<service
  name='network/rsyncd'
  type='service'
  version='1'>
  <create_default_instance enabled='false' />
  <single_instance />

  <dependency name='fs'
    grouping='require_all'
    restart_on='none'
    type='service'>
    <service_fmri value='svc:/system/filesystem/local' />
  </dependency>

  <dependency name='net'
    grouping='require_all'
    restart_on='none'
    type='service'>
    <service_fmri value='svc:/network/loopback' />
  </dependency>

  <exec_method
    type='method'
    name='start'
    exec='/usr/bin/rsync --daemon -v'
    timeout_seconds='5'>
  </exec_method>

  <exec_method
    type='method'
    name='stop'
    exec=':kill'
    timeout_seconds='10'>
  </exec_method>

</service>

</service_bundle>
