# files/galaxy/config/supervisord.conf
[supervisord]
logfile=/var/log/galaxy/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/tmp           ; ('AUTO' child log dir, default $TEMP)
environment = LC_ALL=en_US.utf8,LD_LIBRARY_PATH=/usr/local/galaxy/libpq/


[unix_http_server]
file=/var/tmp/supervisor.sock   	; (the path to the socket file)
chmod=0700                    	; sockef file mode (default 0700)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///var/tmp/supervisor.sock ; use a unix:// URL  for a unix socket



[program:zergpool]
command         = /usr/local/galaxy/galaxy-dist/.venv/bin/uwsgi --master --zerg-pool /usr/local/galaxy/galaxy-dist/var/zergpool.sock:127.0.0.1:4002 --logto /var/log/galaxy/zergpool.log
directory       = /usr/local/galaxy/galaxy-dist/
process_name    = zergpool
priority        = 899
umask           = 022
autostart       = true
autorestart     = true
startsecs       = 5
user            = galaxy
environment     = HOME="/usr/local/galaxy/galaxy-dist/",VIRTUAL_ENV="/usr/local/galaxy/galaxy-dist/.venv",PATH="/usr/local/galaxy/galaxy-dist/.venv/bin:%(ENV_PATH)s"
numprocs        = 1
stopsignal      = INT
stdout_logfile  = /var/log/galaxy/zergpool_supervisor.log
redirect_stderr = true

{% for i in [ 0, 1 ] %}
[program:zergling{{ i }}]
command         = /usr/local/galaxy/galaxy-dist/.venv/bin/uwsgi --ini /usr/local/galaxy/galaxy-dist/config/galaxy.ini --stats 127.0.0.1:919{{ i }} --logto /var/log/galaxy/zergling{{ i }}.log
directory       = /usr/local/galaxy/galaxy-dist/
process_name    = zergling{{ i }}
priority        = 999
umask           = 022
autostart       = {{ "true" if i == 0 else "false"  }}
autorestart     = unexpected
startsecs       = 10
user            = galaxy
environment     = HOME="/usr/local/galaxy/galaxy-dist/",VIRTUAL_ENV="/usr/local/galaxy/galaxy-dist/.venv",PATH="/usr/local/galaxy/galaxy-dist/.venv/bin:%(ENV_PATH)s"
stopsignal      = INT
stopwaitsecs    = 65
stdout_logfile  = /var/log/galaxy/zergling{{ i }}_supervisor.log

{% endfor %}
[program:handler]
command         = /usr/local/galaxy/galaxy-dist/.venv/bin/python /usr/local/galaxy/galaxy-dist/scripts/galaxy-main -c /usr/local/galaxy/galaxy-dist/config/galaxy.ini --server-name=handler%(process_num)s
directory       = /usr/local/galaxy/galaxy-dist/
process_name    = handler%(process_num)s
numprocs        = 12
umask           = 022
autostart       = true
autorestart     = true
startsecs       = 20
user            = galaxy
environment     = VIRTUAL_ENV="/usr/local/galaxy/galaxy-dist/.venv/",PATH="/usr/local/galaxy/galaxy-dist/.venv/bin:%(ENV_PATH)s",GALAXY_WEBAPOLLO_URL="{{ GALAXY_WEBAPOLLO_URL }}",GALAXY_WEBAPOLLO_USER="{{ GALAXY_WEBAPOLLO_USER }}",GALAXY_WEBAPOLLO_PASSWORD="{{ GALAXY_WEBAPOLLO_PASSWORD }}",GALAXY_WEBAPOLLO_EXT_URL="{{ GALAXY_WEBAPOLLO_EXT_URL }}",GALAXY_SHARED_DIR="{{ GALAXY_SHARED_DIR }}"
stdout_logfile  = /var/log/galaxy/handler%(process_num)s.log
stdout_logfile_maxbytes = 128MB
stdout_logfile_backups = 64
redirect_stderr = true



#[program:galaxy_nodejs_proxy]
#directory       = /usr/local/galaxy/galaxy-dist/
#command         = /usr/local/galaxy/galaxy-dist/lib/galaxy/web/proxy/js/lib/main.js --sessions database/session_map.sqlite --ip 127.0.0.1 --port 8800 --reverseProxy
#autostart       = true
#autorestart     = true
#;#user            = galaxy
#startsecs       = 5
#redirect_stderr = true

[program:telegraf_os]
directory       = /usr/local/galaxy/galaxy-dist/
command         = /usr/local/galaxy/galaxy-monit/telegraf --config /usr/local/galaxy/galaxy-monit/telegraf-os.conf
autostart       = true
autorestart     = true
startsecs       = 5
redirect_stderr = true

[program:telegraf_gx]
directory       = /usr/local/galaxy/galaxy-dist/
command         = /usr/local/galaxy/galaxy-monit/telegraf --config /usr/local/galaxy/galaxy-monit/telegraf-gx.conf
autostart       = true
autorestart     = true
startsecs       = 5
redirect_stderr = true

[program:job_conf_service]
directory       = /usr/local/galaxy/galaxy-jobconf/
command         = /usr/local/galaxy/galaxy-jobconf/run.sh
autostart       = true
autorestart     = true
startsecs       = 1
redirect_stderr = true

[program:tiaas_group_join_service]
environment     = CONFIG_PATH="/usr/local/galaxy/tiaas-group-join/config/config.yaml"
directory       = /usr/local/galaxy/tiaas-group-join/code/
command         = /usr/local/galaxy/tiaas-group-join/run.sh
autostart       = true
autorestart     = true
startsecs       = 1
redirect_stderr = true

[group:gx]
programs = zergpool, zergling0, zergling1

[group:hd]
programs = handler
